# ops/server/exambridge-sync.service
[Unit]
Description=Sync ExamBridge from the GitHub production build
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=deploy
Group=deploy
UMask=0027
Environment=EXAMBRIDGE_BASE=/var/www/exambridge
Environment=EXAMBRIDGE_MATERIALS_OWNER=deploy
Environment=EXAMBRIDGE_KEEP_RELEASES=3
Environment=EXAMBRIDGE_HEALTH_URL=https://exambridge.cn/index.html
ExecStart=/var/www/exambridge/shared/sync-gh-pages.sh
TimeoutStartSec=3min
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/www/exambridge
RestrictSUIDSGID=true
LockPersonality=true
CapabilityBoundingSet=
