# ops/server/exambridge-sync.timer
[Unit]
Description=Check GitHub for a new ExamBridge release

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
RandomizedDelaySec=15s
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
